# udev/90-network-udev.rules
# Rule for Network-module to see changing things

ACTION=="add", SUBSYSTEM=="net", RUN+="/usr/local/bin/device_udev add net %p"
ACTION=="add", SUBSYSTEM=="usb", RUN+="/usr/local/bin/device_udev add usb %p"
ACTION=="remove", SUBSYSTEM=="net", RUN+="/usr/local/bin/device_udev del net %p"
ACTION=="remove", SUBSYSTEM=="usb", RUN+="/usr/local/bin/device_udev del usb %p"
